# rootfs-overlay/etc/inittab
::sysinit:/etc/init.d/rcS
null::shutdown:/bin/umount -a -r
